# logitechmediaserver.conf
[supervisord]
nodaemon=true
logfile=/config/supervisord.log

[program:squeezeboxserver]
command=squeezeboxserver --user root --prefsdir /config/prefs --logdir /config/logs --cachedir /config/cache
